# lcd-display.service
[Unit]
Description=LCD Display

[Service]
User=root
Group=root
Restart=on-failure
ExecStart=ksh /root/QnapFreeLCD/lcd-control.ksh

[Install]
WantedBy=multi-user.target
